# infra/ansible/roles/bot/templates/karb-dashboard.service.j2
[Unit]
Description=Karb Dashboard
After=network.target

[Service]
Type=simple
User=karb
Group=karb
WorkingDirectory=/opt/karb
EnvironmentFile=/opt/karb/.env
ExecStart=/opt/karb/venv/bin/karb dashboard --host 0.0.0.0 --port 80
Restart=always
RestartSec=10

# Allow binding to port 80 without root
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=karb-dashboard

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/karb/.karb

[Install]
WantedBy=multi-user.target
